# config/etc/services.d/fluentbit/run
#!/bin/sh

fluent-bit -c /etc//fluent/fluent.conf
